# Rpis/LabPi1/crontab
0 * * * * sudo /usr/local/bin/IPsend
@reboot sudo sleep 10 && sudo /usr/local/bin/IPsend
@reboot /usr/local/bin/remotePDmonitor
@reboot /usr/local/bin/voltreadmonitor
